# bdd-tester-for-osx/services/phantomjs/run
#!/bin/bash

exec /sbin/setuser phantomjs /usr/local/bin/phantomjs --webdriver=4444 > /dev/null
